# systemd/user/ddproulette-roulette.timer
[Unit]
Description=DDP roulette choice timer

[Timer]
OnCalendar=Mon-Fri *-*-* 10:30:00

[Install]
WantedBy=timers.target
